# scripts/supervisord.conf
[supervisord]
nodaemon=true

[program:openclaw-gateway]
command=openclaw gateway --bind lan --port 18789 --allow-unconfigured
autorestart=true

[program:approve-10s]
command=/bin/bash /opt/scripts/approve-10s.sh
autorestart=true